# ansible/roles/google-extrausers-director/templates/google-extrausers-director-sync.service.j2
[Unit]
Description=Sync Google Workspace users -> /var/lib/extrausers (director)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Load SA_KEY, IMPERSONATE, CUSTOMER/DOMAIN, OUTDIR, DB, etc.
EnvironmentFile={{ extrausers_confdir }}/config

# Run via bash so we can expand optional args:
#   - ${DOMAIN:+--domain "$DOMAIN"}  (only added if DOMAIN is set)
#   - ${VERBOSE:+--verbose}          (only added if VERBOSE is non-empty)
ExecStart=/usr/bin/env bash -lc '\
  /usr/bin/python3 {{ extrausers_sbindir }}/google-extrausers-director-sync.py \
    --sa-key "$SA_KEY" \
    --impersonate "$IMPERSONATE" \
    ${CUSTOMER:+--customer "$CUSTOMER"} \
    ${DOMAIN:+--domain "$DOMAIN"} \
    --outdir "${OUTDIR:-/var/lib/extrausers}" \
    --db     "${DB:-/var/lib/googleworkspace-idcache/users.db}" \
    --default-shell "${DEFAULT_SHELL:-/bin/bash}" \
    --home-template "${HOME_TEMPLATE:-/home/{username}}" \
    --group-sync \
    --group-start-gid "${GROUP_START_GID:-30000}" \
    --group-end-gid "${GROUP_END_GID:-39999}" \
    --rps "${RPS:-5}" \
    --max-retries "${MAX_RETRIES:-5}" \
    ${VERBOSE:+--verbose} \
'
# After a successful sync, publish the tarball + sha256 for agents; disabled for now
#ExecStartPost={{ extrausers_sbindir }}/google-extrausers-director-publish

User=root
Group=root
Nice=5

# Hardening (optional but recommended)
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
