# extrausers-pull-rsync/build/lib/systemd/system/extrausers-pull-rsync.service
[Unit]
Description=Pull rendered /var/lib/extrausers from directors
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/extrausers-pull-and-install
User=root
Group=root
Nice=5
